# gunicorn.conf
[program:gunicorn]
command=gunicorn wsgi:app -b localhost:5000
directory=/deploy/MinuteCommute
